# tcp-monitor.service
[Unit]
Description=TCP Session Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
ExecStart=/usr/local/bin/tcp-monitor --config /etc/tcp-monitor/config.toml
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcp-monitor
# Allow the binary to bind privileged ports if needed; remove if ports > 1024.
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
