# tcp-monitor-server.service
[Unit]
Description=TCP Session Monitor (server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nobody
WorkingDirectory=/opt/tcp-monitor
ExecStart=/usr/bin/python3 /opt/tcp-monitor/tcp_monitor_server.py \
    --port 9700 \
    --metrics-port 9701 \
    --heartbeat-interval 30 \
    --log-level INFO
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tcp-monitor-server

[Install]
WantedBy=multi-user.target
